# lib/crontab
0 0 1 * * /etc/lego/renew-certificate.sh 2> /dev/null